# share/toxiproxy.conf
description "TCP proxy to simulate network and system conditions"
author "Simon Eskildsen & Jacob Wirth"

start on startup
stop on shutdown

env HOST="localhost"
env PORT="8474"
env BINARY="/usr/bin/toxiproxy-server"

script
  exec $BINARY -port $PORT -host $HOST
end script
